# systemd/instaguard-engine.service
[Unit]
Description=InstaGuard Firewall Engine
After=network.target nftables.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/instaguard/engine
Environment=INSTAGUARD_DB=/root/instaguard/engine/data/instaguard.db
ExecStart=/root/instaguard/engine/bin/instaguard-engine
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
